# mirror/iOS-shell/Deploy/com.joemoser.dashboard.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  REFERENCE TEMPLATE ONLY. The live LaunchAgent is generated by setup_server.sh
  with machine-specific absolute paths and your Deploy/.env values, and written to
  ~/Library/LaunchAgents/com.joemoser.dashboard.server.plist.
  Placeholders below: @BIN@ = Server/.build/release/Run, @REPO@ = repo root, @LOGS@ = ~/Library/Logs/dashboard.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.joemoser.dashboard.server</string>
  <key>ProgramArguments</key>
  <array>
    <string>@BIN@</string>
    <string>serve</string>
    <string>--env</string>
    <string>production</string>
  </array>
  <key>WorkingDirectory</key><string>@REPO@/Server</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>WEB_DIST</key><string>@REPO@/Web/dist</string>
    <key>PORT</key><string>8080</string>
    <!-- plus TWELVEDATA_KEY, FRED_KEY, DASHBOARD_* injected from Deploy/.env -->
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>@LOGS@/server.out.log</string>
  <key>StandardErrorPath</key><string>@LOGS@/server.err.log</string>
</dict>
</plist>
