# loader/com.joemoser.runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.joemoser.runner</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/joemoser/Dropbox/Source/schedrunner/runner.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <!-- Periodic, exit-quickly job: launchd runs runner.sh every StartInterval
         seconds and lets it exit. KeepAlive must be false so launchd does not
         immediately respawn it (the runner now returns in well under a second). -->
    <key>KeepAlive</key>
    <false/>

    <!-- Do not kill the scripts runner.sh launches when runner.sh exits. This is
         what lets scripts run detached so a slow/hung one never blocks a tick. -->
    <key>AbandonProcessGroup</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/joemoser/Dropbox/Source/schedrunner/runner.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/joemoser/Dropbox/Source/schedrunner/runner.err</string>
  </dict>
</plist>
